# config/supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
command=npm start
autorestart=true
startsecs=0
startretries=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0